# systemd/system/remote-backup.timer
[Unit]
Description=Perform remote backup backups daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
