# systemd/system/remote-backup.timer
[Unit]
Description=Perform remote backup backups daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
